# src/deb/control/postinst
#!/bin/sh

systemctl daemon-reload
echo "Starting Intuso Housemate Globalserver"
systemctl enable housemate-globalserver
systemctl start housemate-globalserver

exit 0